# attendance.service
[Unit]
Description=RFID Attendance Client
After=network.target

[Service]
# Replace 'pi' with your actual username if different
User=admin
WorkingDirectory=/home/admin/rfid-api
ExecStart=/usr/bin/python3 /home/admin/rfid-api/rfid_sender.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
